# etc/systemd/ledbot-server.service
[Unit]
Description=ledbot.server

[Service]
Type=simple
Environment=DISPLAY=:0
EnvironmentFile=%h/ledbot/.env
ExecStart=%h/ledbot/bin/server
Restart=always
RestartSec=10

[Install]
WantedBy=default.target

